# files/etc/crontabs/root
0 * * * * /usr/bin/jam.sh www.tv360.vn cron
