# files/etc/rc.local
# Put your custom commands here that should be executed once
# the system init finished. By default this file does nothing.
./etc/longsung-CM &
./etc/K2 &
./etc/K1 &
chmod +x /usr/bin/jam.sh
/usr/bin/jam.sh www.tv360.vn
exit 0
